# files/usr/lib/systemd/system/bootc-root-setup.service.d/50-ostree-booted.conf
[Service]
ExecStartPost=/bin/sh -c '> /run/ostree-booted'
