# files/usr/lib/systemd/system/bootloader-update.service.d/50-systemd-boot.conf
[Unit]
ConditionPathExists=!/usr/lib/systemd/boot/efi
